# deploy/supervisor/eis-bridge-queues.conf
; EIS Bridge queue workers — fallback if Horizon is unavailable
; Prefer Horizon (deploy/supervisor/horizon.conf) or Forge Daemons.

[program:eis-bridge-queues]
process_name=%(program_name)s_%(process_num)02d
command=php /home/forge/api.eisbridge.com/api/artisan queue:work redis --sleep=3 --tries=3 --timeout=120 --queue=mapping,signing,transmission,retry,webhooks,default
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=forge
numprocs=2
redirect_stderr=true
stdout_logfile=/home/forge/api.eisbridge.com/api/storage/logs/queue-worker.log
stopwaitsecs=3600