# deploy/supervisor/horizon.conf
; EIS Bridge Horizon — Supervisor reference (Forge Daemons or /etc/supervisor/conf.d/)
;
; Forge UI (Site → Daemons) — sandbox.eisbridge.com example:
;   Command:    php8.3 artisan horizon
;   Directory:  /home/forge/sandbox.eisbridge.com/api
;   User:       forge
;   Processes:  1
;
; Production API — replace sandbox path with /home/forge/api.eisbridge.com/api
;
; Forge defaults stopwaitsecs=15; Horizon needs 3600 (Laravel docs). After creating
; the Forge daemon, SSH and set stopwaitsecs=3600 in the generated file under
; /etc/supervisor/conf.d/, then: sudo supervisorctl reread && sudo supervisorctl update
;
; Prerequisites before starting Horizon:
;   - First deploy completed (vendor/, artisan, migrations)
;   - api/.env exists with APP_KEY, QUEUE_CONNECTION=redis, REDIS_* reachable
;   - Redis running on the server
;   - PHP 8.3 with ext-pcntl and ext-posix (Horizon requirement)

[program:eis-bridge-horizon]
process_name=%(program_name)s
command=php8.3 /home/forge/sandbox.eisbridge.com/api/artisan horizon
directory=/home/forge/sandbox.eisbridge.com/api
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=forge
numprocs=1
redirect_stderr=true
stdout_logfile=/home/forge/sandbox.eisbridge.com/api/storage/logs/horizon.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stopwaitsecs=3600