# deploy/supervisor/horizon.conf
; EIS Bridge Horizon — Supervisor reference (Forge Daemons preferred)
; Forge: Site → Daemons → php artisan horizon, directory = .../api
; Adjust paths to match your Forge site (api vs sandbox).

[program:eis-bridge-horizon]
process_name=%(program_name)s
command=php /home/forge/api.eisbridge.com/api/artisan horizon
autostart=true
autorestart=true
user=forge
redirect_stderr=true
stdout_logfile=/home/forge/api.eisbridge.com/api/storage/logs/horizon.log
stopwaitsecs=3600